# packages/xxy/xxy-4574.3.0b0.tar.gz/xxy-4574.3.0b0/deploy/systemd/xxy.service
[Unit]
Description=XXY Chat Bot
Requires=network.target
After=dhcpcd.service
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=simple
User=www-data
Group=www-data
EnvironmentFile=/opt/xxy/.env_systemd
ExecStart=/opt/xxy/venv/bin/python -m uvicorn xxy.web:app --uds /tmp/xxy.sock
WorkingDirectory=/opt/xxy
RestartSec=10
Restart=always
MemoryMax=1000M

[Install]
WantedBy=multi-user.target
